# router/lan/vtluug_iptables.service
[Unit]
Description=VTLUUG iptables
After=network.target

[Service]
ExecStart=/usr/local/bin/vtluug_iptables.sh

[Install]
WantedBy=default.target
